# resticprofile-backup/rootfs/etc/s6-overlay/s6-rc.d/resticprofile/up
/usr/local/bin/resticprofile-schedule.sh
