# affine/rootfs/etc/s6-overlay/s6-rc.d/affine/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

set -euo pipefail

cd /app || exit 1
NODE_ENV="production" node /app/dist/main.js
