# resticprofile-backup/rootfs/etc/s6-overlay/s6-rc.d/crond/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting cron daemon"

crond -f
