# affine/rootfs/etc/s6-overlay/s6-rc.d/affine-migrate/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

set -euo pipefail

bashio::log.info "Running AFFiNE database migrations..."

cd /app || exit 1
NODE_ENV="production" node /app/scripts/self-host-predeploy.js

bashio::log.info "Database migrations completed"
